# departure_times/transitastic.conf
description "Transitastic"
author "chriti@n-bud.de"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid transitastic
setgid transitastic

exec /path/to/run/file.sh